# launchd/io.macvoice.stt.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- 替换 __PYTHON__ 和 __PROJECT_DIR__ 为绝对路径后使用 -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.macvoice.stt</string>
  <key>ProgramArguments</key>
  <array>
    <string>__PYTHON__</string>
    <string>__PROJECT_DIR__/stt_server.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/logs/stt.out.log</string>
  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/logs/stt.err.log</string>
</dict>
</plist>
